# docker_conf/gunicorn.conf
[program:gunicorn]
command=gunicorn  --workers=2 -b 0.0.0.0:5000 texts_viewer.main:app
directory=/deploy/texts_viewer
